# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
pidfile=/var/run/supervisord.pid

[program:ralph-api]
command=/usr/local/bin/ralph-api
# trusted_local = no token needed since nginx proxies from 127.0.0.1
environment=
    RALPH_API_HOST="127.0.0.1",
    RALPH_API_PORT="3000",
    RALPH_API_AUTH_MODE="trusted_local",
    RALPH_API_WORKSPACE_ROOT="/workspace",
    RALPH_API_RALPH_COMMAND="ralph"
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/var/log/supervisor/ralph-api.log
stderr_logfile=/var/log/supervisor/ralph-api.err.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err.log
